# culinary-game-server.service
[Unit]
Description=Culinary Card Game Server
After=network.target

[Service]
Type=simple
User=northpi
WorkingDirectory=/home/northpi/Documents/NorthStar/server
Environment="NODE_ENV=production"
EnvironmentFile=/home/northpi/Documents/NorthStar/server/.env
ExecStart=/usr/bin/node dist/server/src/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

